# ops/systemd/status-metrics.service
[Unit]
Description=Generate status metrics snapshot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/status-metrics.sh
